# apps/macos/macos/LaunchAgents/com.eurora-labs.eurora.launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.eurora-labs.eurora.launcher</string>
	<key>BundleProgram</key>
	<string>Contents/MacOS/EuroraLauncher</string>
	<key>AssociatedBundleIdentifiers</key>
	<array>
		<string>com.eurora-labs.eurora.macos</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<!--
	    Keep the launcher running unconditionally. It is the always-on
	    bridge between Safari and the desktop app; tying its lifecycle
	    to Eurora's would leave the bridge dead between Eurora sessions.
	    The launcher no longer self-terminates on Eurora quit, so this
	    cannot loop.
	-->
	<key>KeepAlive</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
</dict>
</plist>
